# deploy/engine.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.leon.skool-ingest.agent-core</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/leon/Documents/Code/skool-ingest/deploy/run_engine.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/leon/Documents/Code/skool-ingest</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/leon/Documents/Code/skool-ingest/logs/engine.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/leon/Documents/Code/skool-ingest/logs/engine.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
